# config/systemd/user/control-api.service
[Unit]
Description=RustCTA Control API And Dioxus Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/cta/rustcta
Environment=RUST_LOG=info
Environment=RUSTCTA_CONTROL_API_BIND=0.0.0.0:8091
Environment=RUSTCTA_CONTROL_API_ALLOW_NON_CANONICAL_BIND=true
Environment=RUSTCTA_CONTROL_API_AGENT_ID=server-agent
Environment=RUSTCTA_CONTROL_API_TENANT_ID=server
Environment=RUSTCTA_CONTROL_API_AGENT_CAPABILITIES=control-api,supervisor-reader,logs
Environment=RUSTCTA_CONTROL_API_LEGACY_SNAPSHOT_PATH=/home/cta/rustcta/logs/cross_exchange_arbitrage/cross_arb_live_dashboard.json
Environment=RUSTCTA_CONTROL_API_SUPERVISOR_REGISTRY_PATH=/home/cta/rustcta/run/supervisor/registry.json
Environment=RUSTCTA_CONTROL_API_AUDIT_LEDGER_PATH=/home/cta/rustcta/data/control_api/audit.jsonl
Environment=RUSTCTA_CONTROL_API_STRATEGY_LOG_PATH=/home/cta/rustcta/logs/cross_exchange_arbitrage/cross-arb-live.service.log
Environment=RUSTCTA_CONTROL_API_EXCHANGE_API_KEY_STORE=/home/cta/rustcta/data/control_api/exchange_api_keys.env
Environment=RUSTCTA_CONTROL_API_ACCOUNTS_CONFIG=/home/cta/rustcta/config/accounts.yml
Environment=RUSTCTA_CONTROL_API_STATIC_DIR=/home/cta/rustcta/web-ui/dioxus/dist
Environment=RUSTCTA_CONTROL_API_LOCAL_STRATEGY_CONFIG_PATH=/home/cta/rustcta/config/cross_exchange_arbitrage_usdt.yml
Environment=RUSTCTA_CONTROL_API_LOCAL_COMMAND_QUEUE_PATH=/home/cta/rustcta/data/control_api/cross_arb_control_commands.jsonl
Environment=RUSTCTA_CONTROL_API_LOCAL_BALANCE_HISTORY_PATH=/home/cta/rustcta/data/control_api/cross_arb_balance_history.json
Environment=RUSTCTA_CONTROL_API_LOCAL_STRATEGY_PROFIT_HISTORY_PATH=/home/cta/rustcta/logs/cross_exchange_arbitrage/profit_history.jsonl
ExecStart=/home/cta/rustcta/target/release/rustcta-control-api
Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=append:/home/cta/rustcta/logs/control-api.service.log
StandardError=append:/home/cta/rustcta/logs/control-api.service.log

[Install]
WantedBy=default.target
