# linux-deploy/rustcta.service
[Unit]
Description=Rust CTA Trading System
After=network.target
Wants=network.target

[Service]
Type=simple
User=rustcta
Group=rustcta
WorkingDirectory=/opt/rustcta
ExecStart=/opt/rustcta/rustcta-linux
Restart=always
RestartSec=10
KillMode=process
TimeoutStopSec=30

# 环境变量
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rustcta/logs

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustcta

[Install]
WantedBy=multi-user.target